# DataFair.service
[Unit]
Description=DataFair

[Service]
Type=simple
ExecStart=/root/DataFair/DataFair

[Install]
WantedBy=multi-user.target